# os/linux/ubuntu/Shell_Scripts/rc.local
#!/bin/sh -e
sudo mkdir -p /media/redexternal
sudo mount /dev/sdb2 /media/redexternal
exit 0